# infra/systemd/aldermere-api.service
[Unit]
Description=Aldermere API (Hono)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/aldermere/apps/api
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/root/.local/share/fnm/node-versions/v22.23.0/installation/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.local/share/fnm/node-versions/v22.23.0/installation/bin/node /opt/aldermere/apps/api/dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aldermere-api

[Install]
WantedBy=multi-user.target
